# s6-overlay/s6-rc.d/celery-beat/run
#!/command/with-contenv bash
. /app/base.sh
celery -A task beat \
    --loglevel=$LOGLEVEL \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler
