# docker_assets/cron.d/beat-restart
0 * * * * root s6-svc -r /run/service/celery-beat
