# s6-overlay/s6-rc.d/backend/run
#!/command/with-contenv bash
. /etc/s6-overlay/scripts/base.sh
python backend.py
